# packaging/systemd/wordpipe-service.service
[Unit]
Description=Wordpipe Rust dictation service
Documentation=https://github.com/danhansen/wordpipe

[Service]
Type=dbus
BusName=dev.wordpipe.Service
ExecStart=@WORDPIPE_LIBEXEC@/bin/wordpipe-service
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
